# roles/nginx_django/templates/supervisor/django_app.conf.j2
[program:django_app]
command=/django_app/venv/bin/gunicorn --workers 3 --bind unix:/django_app/django_app.sock django_app.wsgi:application
autostart=true
autorestart=true
directory=/django_app
